# utils/supervisord.conf
[unix_http_server]
file=/tmp/supervisor.sock
chmod=0755

[supervisord]
nodaemon=false
logfile=/tmp/supervisord.log
logfile_maxbytes=50MB
logfile_backups=2
loglevel=info
pidfile=/tmp/supervisord.pid

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[program:estserver]
command=estserver -config /docker-estserver/server/server.cfg
directory=/docker-estserver
autostart=true
autorestart=true
startretries=3
stdout_logfile=/tmp/estserver.log
stderr_logfile=/tmp/estserver.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=1